# mailfilter/rootfs/etc/services.d/freshclam/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community App: phpMyAdmin
# Take down the S6 supervision tree when Freshclam fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
    bashio::log.warning "Freshclam crashed, halting app"
    /run/s6/basedir/bin/halt
fi

bashio::log.info "Freshclam stopped, restarting..."