# mailfilter/rootfs/etc/services.d/freshclam/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant App: Mailfilter
# Runs Freshclam virus pattern updates
# ==============================================================================

bashio::log.info "Starting Freshclam...."
exec freshclam -d --config-file /etc/clamav/freshclam.conf >/dev/null 2>&1
