# mailfilter/rootfs/etc/services.d/rspamd/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Mailfilter
# Take down the S6 supervision tree when Rspamd fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
    bashio::log.warning "Rspamd crashed, halting add-on"
    /run/s6/basedir/bin/halt
fi

bashio::log.info "Rspamd stopped, restarting..."